# deployment/scripts/quicklender-worker.service
[Unit]
Description=QuickLender Celery Worker
After=network.target redis.service quicklender-api.service

[Service]
Type=forking
User=quicklender
Group=quicklender
WorkingDirectory=/var/www/quicklender/backend
Environment=DJANGO_ENV=production
EnvironmentFile=/etc/quicklender/.env
ExecStart=/var/www/quicklender/venv/bin/celery     -A quicklender_project worker     --loglevel=info     --concurrency=4     --queues=celery,notifications,mpesa     --logfile=/var/log/quicklender/celery.worker.log     --pidfile=/var/run/quicklender/celery.worker.pid     --detach
ExecStop=/var/www/quicklender/venv/bin/celery     -A quicklender_project control shutdown
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
